# nq-agent.timer
[Unit]
Description=Run nq-agent every three (3) minutes

[Timer]
OnBootSec=15sec
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
